# board/axent/ag11005/skeleton/etc/crontabs/root
# Sync hardware clock.
0 2 * * * /sbin/hwclock -uw
