# board/axent/ag14003/applications/mrwa_vsls/overlay/usr/lib/udev/rules.d/11-sd-cards-auto-mount.rules
KERNEL!="mmcblk[0-9]p[0-9]", GOTO="sd_cards_auto_mount_end"

ACTION=="add", RUN+="/bin/mount /dev/%k"
ACTION=="remove", RUN+="/bin/umount /dev/%k"
LABEL="sd_cards_auto_mount_end"
